# dotroot/etc/systemd/system/cool.service
[Unit]
Description=ISW Cooler Boost

[Service]
Type=simple
Restart=no
ExecStart=/opt/scripts/cool.sh

[Install]
WantedBy=multi-user.target
