# dotroot/etc/systemd/system/cool.service
[Unit]
Description=ISW Cooler Boost
After=suspend.target

[Service]
Type=simple
ExecStart=/usr/bin/isw -b on

[Install]
WantedBy=suspend.target
